# src/ws-ckpt/src/systemd/ws-ckpt.service
# ws-ckpt workspace snapshot daemon
#
# Two ways to start:
#   Service mode:  cp ws-ckpt.service /etc/systemd/system/ && systemctl enable --now ws-ckpt
#   Manual mode:   ws-ckpt daemon [--mount-path ...] [--socket ...] [--log-level debug]

[Unit]
Description=ws-ckpt workspace snapshot daemon
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/ws-ckpt daemon
ExecReload=/usr/bin/ws-ckpt reload
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
